# autobackup.service
[Unit]
Description=ZFS Autobackup Service
Documentation=https://github.com/stein-hak/autobackup
After=network.target zfs-import.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/autobackup
Environment=CONFIG_FILE=/opt/autobackup/backup_config.yaml
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/autobackup/backup_server.py
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal

# Security hardening (optional)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
